# schedule.cron
# attnbot stuff
* * * * * cd /var/apps/smidgeo-headporters && make run-attnbot
