# schedule.cron
# attnbot stuff
8 */4 * * * cd /var/apps/smidgeo-headporters && make run-attnbot-quote
# 38 * * * * cd /var/apps/smidgeo-headporters && make run-attnbot-fact
0 16,19,22 * * * cd /var/apps/smidgeo-headporters && make run-attnbot-popular
0 */8 * * * cd /var/apps/smidgeo-headporters && make run-attnbot-news
0 */6 * * * cd /var/apps/smidgeo-headporters && make run-attnbot-shakespeare
38 */3 * * * cd /var/apps/smidgeo-headporters && make run-attnbot-bible
