# schedule.cron
# attnbot stuff
8 * * * * cd /var/apps/smidgeo-headporters && make run-attnbot
